# home/projedi/.config/systemd/user/mystuff.target
[Unit]
Description=Empty target
After=dbus.service
Requires=dbus.service

[Install]
Alias=default.target
