# home/projedi/.config/systemd/user/timer-daily.timer
[Unit]
Description=Daily timer

[Timer]
OnBootSec=5m
OnUnitActiveSec=1d
Unit=timer-daily.target

[Install]
WantedBy=wm.target
