# home/projedi/.config/systemd/user/urxvtd.service
[Unit]
Description=urxvtd
After=xorg.target

[Service]
ExecStart=/usr/bin/urxvtd -o

[Install]
WantedBy=wm.target
