# etc/systemd/system/conky.service
[Unit]
Description=Conky

[Service]
User=projedi
ExecStart=/usr/bin/conky

[Install]
WantedBy=graphical.target
