# etc/systemd/system/cups.service
[Unit]
Description=CUPS unit

[Service]
ExecStart=/usr/sbin/cupsd -F

[Install]
WantedBy=multi-user.target
