# etc/systemd/system/dropbox@.service
.include /usr/lib/systemd/system/dropbox@.service
[Service]
ExecStartPre=/usr/bin/sleep 30
