# etc/systemd/system/fcron.service
[Unit]
Description=Fcron

[Service]
ExecStart=/usr/sbin/fcron -f

[Install]
WantedBy=multi-user.target
