# etc/systemd/system/mpd.service
[Unit]
Description=Music Player Daemon

[Service]
User=projedi
ExecStart=/usr/bin/mpd --no-daemon

[Install]
WantedBy=multi-user.target
