# etc/systemd/system/ppp@.service
.include /usr/lib/systemd/system/ppp@.service

[Service]
Restart=on-failure
