# etc/systemd/system/sshd.service
[Unit]
Description=SSH Daemon

[Service]
User=projedi
ExecStart=/usr/sbin/sshd -D

[Install]
WantedBy=multi-user.target
