# etc/systemd/system/timer-daily.target
[Unit]
Description=Daily Timer target
StopWhenUnneeded=yes
